# docker/docker-supervisor/nav.ini
[program:nav]
command=/apps/navigator
autostart=true
autorestart=true
stderr_logfile=/var/log/local/nav.err.log
stdout_logfile=/var/log/local/nav.out.log
