# docker/docker-build/bus.ini
[program:bus]
command=/apps/bus
autostart=true
autorestart=true
stderr_logfile=/var/log/local/bus.err.log
stdout_logfile=/var/log/local/bus.out.log
